# .config/systemd/user/jupyterlab.service
[Unit]
Description=Jupyterlab

[Service]
ExecStart=${HOME}/.local/bin/fconda jupyter lab
Type=simple

[Install]
WantedBy=default.target
